# amneziawg-proxy/packaging/amneziawg-proxy.service
[Unit]
Description=AmneziaWG UDP Proxy (DPI obfuscation layer)
After=network.target
Documentation=https://github.com/wiresock/amneziawg-install

[Service]
Type=simple
Restart=on-failure
RestartSec=5s

# --- Run as root to bind privileged / low-numbered ports ---------------------
# The proxy may need to bind to ports below 1024 (e.g. 443/UDP for QUIC).
# If the listen port is ≥ 1024 you can drop root by adding:
#   User=nobody
#   AmbientCapabilities=CAP_NET_BIND_SERVICE
User=root

# --- Binary and configuration ------------------------------------------------
ExecStart=/usr/local/bin/amneziawg-proxy /etc/amneziawg-proxy/proxy.toml
WorkingDirectory=/var/lib/amneziawg-proxy

# --- Hardening ---------------------------------------------------------------
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Allow reading the AWG server config and the proxy config.
ReadOnlyPaths=/etc/amnezia /etc/amneziawg-proxy

# Allow writing to the working directory (reserved for future state files).
ReadWritePaths=/var/lib/amneziawg-proxy

# --- Logging -----------------------------------------------------------------
# Adjust RUST_LOG for more verbose output: RUST_LOG=amneziawg_proxy=debug
Environment=RUST_LOG=amneziawg_proxy=info

[Install]
WantedBy=multi-user.target
